# archiso/releng/ai-assist-launcher.service
[Unit]
Description=KENKI AI Assistant Launcher
After=network.target

[Service]
Type=simple
ExecStart=/usr/bin/python3 /etc/kenki/ai-assist/kenki_assist.py --interactive
Restart=on-failure

[Install]
WantedBy=default.target 